# docker-python/delivery-service/Docker/crontab/crontab
#run python script every minutes
* * * * * /usr/local/bin/python /app/main.py >> /app/log.txt
* * * * * sleep 30; /usr/local/bin/python /app/main.py >> /app/log.txt
